# hosts/masternode/systemd/vmstation-idle-check.service
[Unit]
Description=VMStation idle check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vmstation_idle_check.sh
